# com.user.autoraise-rs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.autoraise-rs</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Update this path to wherever you install the binary -->
        <string>/usr/local/bin/autoraise-rs</string>
        <!-- Optional flags — remove or adjust as needed -->
        <string>--aerospace-aware</string>
        <string>true</string>
        <string>--delay</string>
        <string>1</string>
        <string>--poll-millis</string>
        <string>50</string>
        <string>--disable-key</string>
        <string>control</string>
    </array>

    <!-- Run at login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart if it crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Log output (create /var/log manually or change path) -->
    <key>StandardOutPath</key>
    <string>/tmp/autoraise-rs.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/autoraise-rs.err</string>

    <!-- Low priority — don't compete with your actual work -->
    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
